# debian/powerwaked.powerwaked.service
[Unit]
Description=PowerWake Server

[Service]
Type=simple
ExecStart=/usr/sbin/powerwaked

[Install]
WantedBy=multi-user.target
